# systemd/dyfi.timer
[Unit]
Description=dy.fi domain update timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1d
Unit=dyfi.service

[Install]
WantedBy=network.target
